# docker/hadoop/namenode/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:pythonapp]
command=/run.sh

[program:sparkhistoryserver]
command=/run_spark_history_server.sh